# supervisord.conf
[supervisord]
nodaemon=false
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:dashboard]
command=python dashboard/app.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/dashboard.log
stderr_logfile=/app/logs/dashboard.log
priority=10

[program:paper_breakout]
command=python -u paper/paper_trade_breakout.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/paper_breakout.log
stderr_logfile=/app/logs/paper_breakout.log
priority=20

[program:paper_calmar]
command=python -u paper/paper_trade_calmar.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/paper_calmar.log
stderr_logfile=/app/logs/paper_calmar.log
priority=20

[program:paper_regime]
command=python -u paper/paper_trade_regime.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/paper_regime.log
stderr_logfile=/app/logs/paper_regime.log
priority=20

[program:paper_martingale]
command=python -u paper/paper_trade_martingale.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/paper_martingale.log
stderr_logfile=/app/logs/paper_martingale.log
priority=20

[program:paper_boll_scalp]
command=python -u paper/paper_trade_boll_scalp.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/paper_boll_scalp.log
stderr_logfile=/app/logs/paper_boll_scalp.log
priority=20
